# containers/django/config/s6/gunicorn/run
#!/command/execlineb -P
with-contenv
export PYTHONPATH "/app"
/bin/uv
run
/app/.venv/bin/python3
-m
gunicorn
soltyback.wsgi:application
--bind
0.0.0.0:8000